# 99-backlight.rules
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="rpi_backlight", RUN+="/bin/chmod a+w /sys/class/backlight/%k/brightness"
